# docker/root/etc/services.d/postfix/run
#!/usr/bin/with-contenv bash

echo "Starting Mail Server"

# Update aliases
/usr/sbin/postalias /etc/aliases

exec \
	/usr/sbin/postfix start-fg